# @kelp/services/worker/crontab
*/2 * * * * createPoe
*/1 * * * * mediaPoeExists
*/1 * * * * devicePoeExists
